# selenium-patches/browserAutomation.conf
[program:browserAutomation]
command=/opt/bin/browserAutomation
autostart=true
autorestart=true
user=${SEL_UID}
