# selenium-patches/browserAutomation.conf
[program:browserAutomation]
command=/opt/bin/rbee
autostart=true
autorestart=true
user=seluser
